# mtpromo.conf
[program:mtpromo]
command=/usr/bin/python3 /home/mtpromo/player.py
autostart=true
autorestart=true
